# docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --queue=default --daemon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=4
redirect_stderr=true
stdout_logfile=/src/storage/logs/worker.log
stderr_logfile=/src/storage/logs/worker-error.log
stopwaitsecs=3600
